# etc/profile.d/emsdk.sh
export PATH="$HOME/src/github.com/emscripten-core/emsdk:$PATH"

export EMSDK_QUIET=1

[ -f $HOME/src/github.com/emscripten-core/emsdk/emsdk_env.sh ] && source $HOME/src/github.com/emscripten-core/emsdk/emsdk_env.sh
