# etc/profile.d/homebrew.sh
if command -v brew &> /dev/null ; then
  export PATH="$(brew --prefix):$PATH"
  export PKG_CONFIG_PATH="$(brew --prefix)/lib/pkgconfig:$(brew --prefix)/lib:$PKG_CONFIG_PATH"

  [[ -r "$(brew --prefix)/etc/profile.d/bash_completion.sh" ]] && . "$(brew --prefix)/etc/profile.d/bash_completion.sh"

  export HOMEBREW_NO_ANALYTICS=1
  export HOMEBREW_NO_AUTO_UPDATE=1
fi

